# godns/rootfs/etc/s6-overlay/s6-rc.d/godns/finish
#!/usr/bin/env bashio

if test "$1" -eq 256 ; then
  bashio::log.warning "Halt add-on with exit code ${1}"
  e=$((128 + $2))
else
  e="$1"
fi

echo "$e" > /run/s6-linux-init-container-results/exitcode
# exec /run/s6/basedir/bin/halt

bashio::log.info "Service restart after closing"